# scripts/kestrel-aspnetcoreservice.service
[Unit]
Description=AWS CodeStar Sample ASP.NET Core Web Service

[Service]
WorkingDirectory=/home/ubuntu/migratedVisitorCounter/
ExecStart=/usr/bin/dotnet /home/ubuntu/migratedVisitorCounter/MigratedVisitorCounter.dll
Restart=always
# Restart service after 10 seconds if dotnet service crashes
RestartSec=10
SyslogIdentifier=migratedVisitorCounter
User=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production 

[Install]
WantedBy=multi-user.target